# prometheus-pbs-exporter.timer
[Unit]
Description=Prometheus PBS Exporter Timer

[Timer]
OnCalendar=*-*-* 08:00:00

[Install]
WantedBy=timers.target
